# systemd/stratusd.service
[Unit]
Description=Stratus VM Orchestrator Daemon
Documentation=https://github.com/sdmitriev1/stratus
After=network.target
Requires=stratusd.socket

[Service]
Type=notify
ExecStart=/usr/local/bin/stratusd
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/lib/stratus /run/stratus /sys/fs/bpf/stratus
ProtectHome=yes

# Required capabilities
AmbientCapabilities=CAP_NET_ADMIN CAP_BPF CAP_SYS_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_BPF CAP_SYS_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
